# debian/profiles/debian/jessie/oasis-profile-as-system.postinst
#!/bin/bash

# Enable and start asterisk
/bin/systemctl enable asterisk
/bin/systemctl start asterisk

:
