# debian/users/debian/jessie/postinst
#!/bin/bash

## This package contains systemd unit files
/bin/systemctl daemon-reload
/bin/systemctl enable kamailio-proxyusers

:
