# debian/systemd/cdrparser.timer
[Unit]
Description=Run CDRparser timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
