# debian/systemd/cdrparser.service
[Unit]
Description=CDRparser

[Service]
Type=oneshot
ExecStart=/opt/irontec/ivozprovider/scripts/CDRparser.pl
