# debian/systemd/fastagi.socket
[Socket]
ListenStream=127.0.0.1:4573
Accept=true

[Install]
WantedBy=sockets.target
