# debian/systemd/ivozprovider-scheduler.service
[Unit]
Description=Run Ivozprovider Scheduled tasks

[Service]
Type=oneshot
ExecStart=/opt/irontec/ivozprovider/microservices/scheduler/bin/run
