# debian/systemd/ivozprovider-balances.timer
[Unit]
Description=Run Ivozprovider Balances updater

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
