# debian/systemd/ivozprovider-jwt.timer
[Unit]
Description=Ivozprovider JWT tasks runner

[Timer]
OnBootSec=1h
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
