# debian/systemd/ivozprovider-realtime.service
[Unit]
Description=Ivozprovider Realtime Websockets Server

[Service]
User=root
ExecStart=/opt/irontec/ivozprovider/microservices/realtime/bin/run
