# debian/systemd/ivozprovider-users-cdrs.timer
[Unit]
Description=Ivozprovider Users CDR mixer

[Timer]
OnBootSec=3min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
