# debian/systemd/ivozprovider-recordings.timer
[Unit]
Description=Run Ivozprovider Recordings handler

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
