# microservices/realtime/ivozprovider-realtime.service
[Unit]
Description=Ivozprovider Realtime Event publisher

[Service]
User=root
ExecStart=/opt/irontec/ivozprovider/microservices/realtime/bin/run
