# debian/ivozprovider-asterisk-agi.fastagi.socket
[Socket]
ListenStream=127.0.0.1:4573

[Install]
WantedBy=sockets.target
