# debian/ivozprovider-profile-data.postinst
#!/bin/bash

# Enable and start mysql
/bin/systemctl enable mysql
/bin/systemctl start mysql

# Enable and start cdr-parser
/bin/systemctl daemon-reload
/bin/systemctl start cdrparser.timer
/bin/systemctl enable cdrparser.timer

:

