# profiles/portal/etc/supervisor/conf.d/cgrates-daemon.conf
[program:cgrates]
command=/usr/bin/php /opt/irontec/ivozprovider/microservices/workers/bin/cgrates-reloader
autorestart=true
autostart=true
user=www-data
startretries=5
process_name=%(program_name)s-%(process_num)s
numprocs=1
stopsignal=KILL
