# profiles/portal/etc/supervisor/conf.d/multimedia-daemon.conf
[program:multimedia]
command=/usr/bin/php /opt/irontec/ivozprovider/microservices/workers/bin/multimedia-encoder
autorestart=true
autostart=true
user=www-data
startretries=5
process_name=%(program_name)s-%(process_num)s
numprocs=1
stopsignal=KILL
