# profiles/portal/etc/supervisor/conf.d/rates-importer.conf
[program:ratesimporter]
command=/usr/bin/php /opt/irontec/ivozprovider/microservices/workers/bin/rates-import
autorestart=true
autostart=true
user=www-data
startretries=5
process_name=%(program_name)s-%(process_num)s
numprocs=1
stopsignal=KILL
