# profiles/portal/etc/supervisor/conf.d/xmlrpcdelayedproxytrunks-daemon.conf
[program:xmlrpcdelayedproxytrunks]
command=/usr/bin/php /opt/irontec/ivozprovider/portals/application/bin/cli.php -a default/worker/xmlrpcdelayedproxytrunks -e %(ENV_APPLICATION_ENV)s
autorestart=true
autostart=true
user=www-data
startretries=100000000

 