# profiles/data/lib/systemd/system/CDRparser.service
[Unit]
Description=CDRparser

[Service]
Type=oneshot
ExecStart=/opt/irontec/ivozprovider/cdr-parser/CDRparser.pl
