# profiles/data-system/lib/systemd/system/CDRparser.timer
[Unit]
Description=Run CDRparser every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
