# profiles/data-system/lib/systemd/system/CDRparser.service
[Unit]
Description=CDRparser

[Service]
Type=oneshot
ExecStart=/opt/oasis/cdr-parser/CDRparser.pl
